# scripts/mark-l-agent-hardened.service
# Mark-L Agent — systemd service（沙箱版）
#
# 适用于仅 Chat 模式（不启用 Dev Agent 的 Bash 工具）。
# ProtectSystem=strict 会限制文件系统写入，需用 ReadWritePaths 显式放行。
#
# 安装方式同 mark-l-agent.service。
# 按实际路径修改 WorkingDirectory、ExecStart、ReadWritePaths。

[Unit]
Description=Mark-L Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mark-l-agent
Group=mark-l-agent
WorkingDirectory=/home/mark-l-agent/mark-l-agent
ExecStart=/home/mark-l-agent/mark-l-agent/.venv/bin/python main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
BindPaths=/home/mark-l-agent/mark-l-agent
ReadWritePaths=/home/mark-l-agent/mark-l-agent/data
ReadWritePaths=/home/mark-l-agent/mark-l-agent/audit.log
ReadWritePaths=/home/mark-l-agent/mark-l-agent/permissions.json

[Install]
WantedBy=multi-user.target
